# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:webpage-content-extractor]
command=/webpage-content-extractor/bin/www
autorestart=true
